# packages/011-telemetry-storage-rollups/011-telemetry-storage-rollups/systemd/nexus-telemetry.timer
[Unit]
Description=Collect Nexus Platform telemetry every minute
[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
AccuracySec=10s
Persistent=true
Unit=nexus-telemetry.service
[Install]
WantedBy=timers.target
